# deploy/com.llw.feishu-daily-work.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.llw.feishu-daily-work</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/ccrt/Library/Application Support/LLW Assistant/components/feishu-daily-work/src/main.mjs</string>
    <string>/Users/ccrt/Library/Application Support/LLW Assistant/state/feishu-daily-work/config.json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ccrt/Library/Application Support/LLW Assistant/components/feishu-daily-work</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict><key>SuccessfulExit</key><false/></dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/Users/ccrt/Library/Application Support/LLW Assistant/logs/feishu-daily-work/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ccrt/Library/Application Support/LLW Assistant/logs/feishu-daily-work/stderr.log</string>
</dict>
</plist>
